# phosphor-hostEvent-monitor.service.in
[Unit]
Description=BMC hostEvent Monitor
After=xyz.openbmc_project.Chassis.Control.Power.service

[Service]
ExecStart=@bindir@/hostEvent-monitor
Type=dbus
BusName=xyz.openbmc_project.HostEventMon
SyslogIdentifier=phosphor-hostEvent-monitor

[Install]
WantedBy=multi-user.target
